# 02_recursos/operativos.socket
[Unit]
Description=Socket de la API SGI Operativos

[Socket]
ListenStream=/run/operativos.sock

[Install]
WantedBy=sockets.target
